# systemd/raspberry/bdx-environment-ioc.service.in
[Unit]
Description=BDX Raspberry MCP9808 environment IOC
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=@BDX_RUNTIME_USER@
Group=@BDX_RUNTIME_GROUP@
SupplementaryGroups=i2c
WorkingDirectory=/opt/bdx-slow-control
EnvironmentFile=/etc/bdx-slow-control/bdx.env
ExecStart=/opt/bdx-slow-control/.venv/bin/bdx-environment-ioc --config /etc/bdx-slow-control/profiles/raspberry/environment.json
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
