# systemd/main-server/bdx-main-server-ioc.service.in
[Unit]
Description=BDX main-server IOC
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=@BDX_RUNTIME_USER@
Group=@BDX_RUNTIME_GROUP@
WorkingDirectory=/opt/bdx-slow-control
EnvironmentFile=/etc/bdx-slow-control/bdx.env
ExecStart=/opt/bdx-slow-control/.venv/bin/bdx-prototype-ioc --config-dir /etc/bdx-slow-control/profiles/main-server
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
