# systemd/split-host/bdx-daq-ioc.service
[Unit]
Description=BDX DAQ interface IOC
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=streamdaq
Group=streamdaq
WorkingDirectory=/opt/bdx-slow-control
EnvironmentFile=-/etc/bdx-slow-control/bdx.env
ExecStart=/opt/bdx-slow-control/.venv/bin/bdx-daq-ioc --config /etc/bdx-slow-control/daq.json
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
